# images/cuda/ROOT/etc/s6-overlay/s6-rc.d/code-server/run
#!/bin/sh
set -e;

export HOME=/home/user

. /venv/main/bin/activate

exec code-server --disable-telemetry --auth none /workspace
