# images/cuda/ROOT-aitoolkit/etc/s6-overlay/s6-rc.d/aitoolkit/run
#!/command/with-contenv sh
set -e

case "${AITOOLKIT_ENABLED:-1}" in
    0|false|False|FALSE|no|No|NO|off|Off|OFF)
        echo "AI Toolkit disabled via AITOOLKIT_ENABLED=${AITOOLKIT_ENABLED}"
        exec /usr/bin/sleep infinity
        ;;
esac

. /venv/aitoolkit/bin/activate
cd /opt/workspace-internal/ai-toolkit/ui
exec npm run start -- --port ${AITOOLKIT_PORT:-8675}
