# images/cuda/ROOT/etc/s6-overlay/s6-rc.d/comfy/run
#!/bin/bash
set -e;

source /venv/main/bin/activate

cd /opt/workspace-internal/ComfyUI/
if [ -d /workspace/basedir ]; then
    echo "Directory /workspace/basedir exists."
else
    echo "Directory /workspace/basedir does not exist. Creating it."
    mkdir -p /workspace/basedir
    cp -a /opt/workspace-internal/ComfyUI/custom_nodes/ /workspace/basedir/
fi

export LD_PRELOAD=libtcmalloc_minimal.so.4

# Reduce allocator fragmentation unless explicitly overridden.
export PYTORCH_CUDA_ALLOC_CONF="${PYTORCH_CUDA_ALLOC_CONF:-expandable_segments:True,max_split_size_mb:128}"

# Default to no VRAM mode for better compatibility on <=10GB GPUs.
# Override with COMFYUI_VRAM_FLAG (e.g. --normalvram, --lowvram, --novram).
vram_flag="${COMFYUI_VRAM_FLAG:---novram}"

exec python main.py --preview-method auto --base-directory /workspace/basedir ${vram_flag} ${COMFYUI_ARGS}