# images/cuda/ROOT/etc/s6-overlay/s6-rc.d/comfy/run
#!/command/with-contenv sh
set -e

. /venv/main/bin/activate

case "${COMFYUI_ENABLED:-1}" in
    0|false|False|FALSE|no|No|NO|off|Off|OFF)
        echo "ComfyUI disabled via COMFYUI_ENABLED=${COMFYUI_ENABLED}"
    exec /usr/bin/sleep infinity
        ;;
esac

cd /opt/workspace-internal/ComfyUI/
if [ -d /workspace/basedir ]; then
    echo "Directory /workspace/basedir exists."
else
    echo "Directory /workspace/basedir does not exist. Creating it."
    mkdir -p /workspace/basedir
    cp -a /opt/workspace-internal/ComfyUI/custom_nodes/ /workspace/basedir/
fi

export LD_PRELOAD=libtcmalloc_minimal.so.4

exec python main.py --base-directory /workspace/basedir ${COMFYUI_VRAM_FLAG} ${COMFYUI_ARGS}