# images/cuda/ROOT/etc/s6-overlay/s6-rc.d/comfy/run
#!/command/with-contenv sh
set -e

. /venv/main/bin/activate

case "${COMFYUI_ENABLED:-1}" in
    0|false|False|FALSE|no|No|NO|off|Off|OFF)
        echo "ComfyUI disabled via COMFYUI_ENABLED=${COMFYUI_ENABLED}"
        COMFYUI_START_FLAG="${COMFYUI_START_FLAG:-/workspace/.start-comfyui}"
        echo "Waiting for start flag: ${COMFYUI_START_FLAG}"
        while [ ! -f "${COMFYUI_START_FLAG}" ]; do
            /usr/bin/sleep 5
        done
        echo "Start flag detected, launching ComfyUI"
        ;;
esac

cd /opt/workspace-internal/ComfyUI/

mkdir -p /workspace/basedir
if [ ! -d /workspace/basedir/custom_nodes ]; then
    echo "custom_nodes directory does not exist in /workspace/basedir. Copying from /opt/workspace-internal/ComfyUI/custom_nodes."
    cp -a /opt/workspace-internal/ComfyUI/custom_nodes/ /workspace/basedir/
fi

export LD_PRELOAD=libtcmalloc_minimal.so.4

if [ "$SAGEATTENTION_USE" != "0" ] && { [ "$SAGEATTENTION_VERSION" = "v2" ] || [ "$SAGEATTENTION_VERSION" = "v3" ]; }; then
    exec python main.py --base-directory /workspace/basedir --use-sage-attention ${COMFYUI_VRAM_FLAG} ${COMFYUI_ARGS}
else
    exec python main.py --base-directory /workspace/basedir ${COMFYUI_VRAM_FLAG} ${COMFYUI_ARGS}
fi
