# images/cuda/ROOT/etc/s6-overlay/s6-rc.d/caddy/run
#!/bin/bash
set -e;

cd /etc/caddy;
/command/with-contenv gencert;
exec caddy run